# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log

[program:main_app]
command=python app.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/main_app.log
stderr_logfile=/var/log/supervisor/main_app.log
environment=HF_SPACE="1",GRADIO_SERVER_PORT="7862",GRADIO_SERVER_NAME="0.0.0.0"

[program:preview_app]
command=python sandbox/app.py --server-port 7861 --server-name 0.0.0.0
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/preview_app.log
stderr_logfile=/var/log/supervisor/preview_app.log
